# systemd/.config/systemd/user/waybar.service
[Unit]
Description=Waybar Service

[Service]
ExecStart=/usr/bin/waybar
Restart=always

[Install]
WantedBy=graphical-session.target